# local-run/build-image/runit/salt-master/run
#!/bin/sh
if [ $SALT_MASTER_ENABLE -eq 1 ]; then
salt-master
else
sleep infinity
fi
